# legacy/supervisord.conf
; Supervisord layout for production: one program per Celery queue family + a
; singleton beat. Workers are split so a long-running parse / agent batch
; never starves the fast queues (webhooks, dispatcher tick).
;
; Horizontal scale: deploy this image as N replicas behind a load balancer.
; Replicas are interchangeable EXCEPT beat must run in exactly one container
; — set ``BEAT_ENABLED=false`` on every replica except the leader. The
; dispatcher task itself is idempotent (``FOR UPDATE SKIP LOCKED``) so
; accidental double-fires are safe, but two beats compounding cron drift is
; not what you want.
;
; Per-program concurrency is env-tunable so each replica can match its CPU /
; memory budget without rebuilding the image.
;
; Note on supervisord substitution: ``%`` starts a Python-style substitution,
; so a literal ``%`` (e.g. for ``%h`` hostname expansion in ``celery -n``)
; must be doubled to ``%%``.

[supervisord]
nodaemon=true
loglevel=info

; ─────────────────────────────────────────────────────────────────────────
; API server
; ─────────────────────────────────────────────────────────────────────────
[program:gunicorn]
command=/bin/bash -c 'echo "Starting Gunicorn..." && if [ -f /app/newrelic.ini ]; then exec uv run newrelic-admin run-program gunicorn --workers $(nproc) --worker-class uvicorn.workers.UvicornWorker --timeout 1800 --bind 0.0.0.0:8001 --log-level info app.main:app; else exec uv run gunicorn --workers $(nproc) --worker-class uvicorn.workers.UvicornWorker --timeout 1800 --bind 0.0.0.0:8001 --log-level info app.main:app; fi'
autostart=true
autorestart=true
priority=10
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=root

; ─────────────────────────────────────────────────────────────────────────
; Celery: parsing (process_repository) — heavy, memory-hungry, long-running
; ─────────────────────────────────────────────────────────────────────────
[program:celery_parsing]
command=/bin/bash -c 'echo "Starting Celery: parsing..." && if [ -f /app/newrelic.ini ]; then exec uv run newrelic-admin run-program celery -A app.celery.celery_app worker --loglevel=info -Q ${CELERY_QUEUE_NAME}_process_repository -n parsing@%%h -E --concurrency=${CELERY_PARSING_CONCURRENCY:-2}; else exec uv run celery -A app.celery.celery_app worker --loglevel=info -Q ${CELERY_QUEUE_NAME}_process_repository -n parsing@%%h -E --concurrency=${CELERY_PARSING_CONCURRENCY:-2}; fi'
autostart=true
autorestart=true
priority=20
stopwaitsecs=600
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=root

; ─────────────────────────────────────────────────────────────────────────
; Celery: agent tasks — LLM calls, mid-duration, IO-bound
; ─────────────────────────────────────────────────────────────────────────
[program:celery_agent]
command=/bin/bash -c 'echo "Starting Celery: agent..." && if [ -f /app/newrelic.ini ]; then exec uv run newrelic-admin run-program celery -A app.celery.celery_app worker --loglevel=info -Q ${CELERY_QUEUE_NAME}_agent_tasks -n agent@%%h -E --concurrency=${CELERY_AGENT_CONCURRENCY:-4}; else exec uv run celery -A app.celery.celery_app worker --loglevel=info -Q ${CELERY_QUEUE_NAME}_agent_tasks -n agent@%%h -E --concurrency=${CELERY_AGENT_CONCURRENCY:-4}; fi'
autostart=true
autorestart=true
priority=20
stopwaitsecs=600
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=root

; ─────────────────────────────────────────────────────────────────────────
; Celery: webhooks / external events — fast, bursty
; ─────────────────────────────────────────────────────────────────────────
[program:celery_events]
command=/bin/bash -c 'echo "Starting Celery: events..." && if [ -f /app/newrelic.ini ]; then exec uv run newrelic-admin run-program celery -A app.celery.celery_app worker --loglevel=info -Q external-event -n events@%%h -E --concurrency=${CELERY_EVENTS_CONCURRENCY:-4}; else exec uv run celery -A app.celery.celery_app worker --loglevel=info -Q external-event -n events@%%h -E --concurrency=${CELERY_EVENTS_CONCURRENCY:-4}; fi'
autostart=true
autorestart=true
priority=20
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=root

; ─────────────────────────────────────────────────────────────────────────
; Celery: context-graph reconciliation (dispatcher tick + agent batches).
; Skips startup when CONTEXT_GRAPH_ENABLED is false/0/no/off so unconsumed
; messages don't pile up in Redis. Multiple replicas can each run this and
; FOR UPDATE SKIP LOCKED on the batch claim handles concurrent draining.
; ─────────────────────────────────────────────────────────────────────────
[program:celery_context]
command=/bin/bash -c 'cg="${CONTEXT_GRAPH_ENABLED:-true}"; cg_lc=$(printf "%%s" "$cg" | tr "[:upper:]" "[:lower:]"); if [ "$cg_lc" = "false" ] || [ "$cg_lc" = "0" ] || [ "$cg_lc" = "no" ] || [ "$cg_lc" = "off" ]; then echo "CONTEXT_GRAPH_ENABLED=$cg → context-graph worker not starting"; exit 0; fi; echo "Starting Celery: context-graph..." && if [ -f /app/newrelic.ini ]; then exec uv run newrelic-admin run-program celery -A app.celery.celery_app worker --loglevel=info -Q context-graph-etl -n context@%%h -E --concurrency=${CELERY_CONTEXT_CONCURRENCY:-2}; else exec uv run celery -A app.celery.celery_app worker --loglevel=info -Q context-graph-etl -n context@%%h -E --concurrency=${CELERY_CONTEXT_CONCURRENCY:-2}; fi'
autostart=true
autorestart=unexpected
exitcodes=0
priority=20
stopwaitsecs=600
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=root

; ─────────────────────────────────────────────────────────────────────────
; Celery beat — cluster-wide singleton.
; Set ``BEAT_ENABLED=false`` on every replica EXCEPT one. The single beat
; instance publishes scheduled messages (e.g. the dispatcher tick) to the
; broker; any worker on the relevant queue picks them up.
;
; ``--pidfile=`` (empty) tells beat not to write a pidfile, which avoids
; stale-lock issues across supervisord restarts.
; ─────────────────────────────────────────────────────────────────────────
[program:celery_beat]
command=/bin/bash -c 'be="${BEAT_ENABLED:-true}"; be_lc=$(printf "%%s" "$be" | tr "[:upper:]" "[:lower:]"); if [ "$be_lc" = "false" ] || [ "$be_lc" = "0" ] || [ "$be_lc" = "no" ] || [ "$be_lc" = "off" ]; then echo "BEAT_ENABLED=$be → celery beat not starting"; exit 0; fi; echo "Starting Celery beat..." && if [ -f /app/newrelic.ini ]; then exec uv run newrelic-admin run-program celery -A app.celery.celery_app beat --loglevel=info --schedule=/tmp/celerybeat-schedule --pidfile=; else exec uv run celery -A app.celery.celery_app beat --loglevel=info --schedule=/tmp/celerybeat-schedule --pidfile=; fi'
autostart=true
autorestart=unexpected
exitcodes=0
priority=30
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=root
